# server/yt2epub-brief.service
[Unit]
Description=yt2epub daily brief (one-shot, triggered by timer)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=yt2epub
WorkingDirectory=/home/yt2epub/yt2epub
ExecStart=/usr/bin/python3 -u /home/yt2epub/yt2epub/daily_brief.py
Environment="PYTHONUNBUFFERED=1"
StandardOutput=append:/home/yt2epub/yt2epub/daily_brief.stdout.log
StandardError=append:/home/yt2epub/yt2epub/daily_brief.stderr.log
