# server/yt2epub-bot.service
[Unit]
Description=yt2epub Telegram bot service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yt2epub
WorkingDirectory=/home/yt2epub/yt2epub
ExecStart=/usr/bin/python3 -u /home/yt2epub/yt2epub/bot_service.py
Restart=on-failure
RestartSec=30
Environment="PYTHONUNBUFFERED=1"
StandardOutput=append:/home/yt2epub/yt2epub/bot_service.stdout.log
StandardError=append:/home/yt2epub/yt2epub/bot_service.stderr.log

[Install]
WantedBy=multi-user.target
